# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:postgres]
command=/usr/lib/postgresql/16/bin/postgres -D /var/lib/postgresql/16/main -c config_file=/etc/postgresql/16/main/postgresql.conf
user=postgres
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/postgres-stdout.log
stderr_logfile=/var/log/supervisor/postgres-stderr.log
environment=PGDATA=/var/lib/postgresql/16/main,TZ="Africa/Cairo"

[program:backend]
command=uvicorn main:app --host 0.0.0.0 --port 8000
directory=/app
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/backend-stdout.log
stderr_logfile=/var/log/supervisor/backend-stderr.log
environment=POSTGRES_HOST="localhost",POSTGRES_PORT="5432",POSTGRES_DB="eviltwin",POSTGRES_USER="eviltwin",TZ="Africa/Cairo"

[program:nginx]
command=nginx -g "daemon off;"
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx-stdout.log
stderr_logfile=/var/log/supervisor/nginx-stderr.log
